# 80-keychron.rules
SUBSYSTEMS=="usb", ATTRS{product}=="Keychron K3", RUN+="/bin/sh -c 'echo 2 | tee /sys/module/hid_apple/parameters/fnmode'"
